# local/etc/profile.d/colorsvn.sh
# set alias for svn

alias svn=colorsvn
