# local/etc/profile.d/gnupg.sh
# See gpg-agent(1).
export GPG_TTY=`tty`
