# local/etc/profile.d/vagrant.sh
if [ -e ~/.vagrant.d/license-vagrant-vmware-fusion.lic ]; then
  export VAGRANT_DEFAULT_PROVIDER=vmware_fusion
fi
